# services/arrpc.service
[Unit]
Description=arRPC Service for Vesktop

[Service]
Type=simple
ExecStart=arrpc
